# Tick_Data_Folder/tick_harvester.service
# ─────────────────────────────────────────────────────────────────────────────
#  tick_harvester.service
#  systemd unit file for the Sharekhan Tick Harvester
#  Target: AWS EC2 t3.small  (2 GB RAM, 2 vCPU)  —  Ubuntu 22.04 LTS
#
# ── ONE-TIME EC2 SETUP ───────────────────────────────────────────────────────
#
#  1.  Set the instance timezone to IST (do this once after launch):
#          sudo timedatectl set-timezone Asia/Kolkata
#          timedatectl    # verify: "Time zone: Asia/Kolkata (IST, +0530)"
#
#  2.  Clone / upload the project:
#          scp -r "Stock-Predictor-v1.1.0" ubuntu@<EC2-IP>:/home/ubuntu/
#      Or clone from git:
#          git clone <repo-url> /home/ubuntu/Stock-Predictor
#
#  3.  Create a Python virtualenv and install dependencies:
#          cd /home/ubuntu/Stock-Predictor
#          python3 -m venv venv
#          venv/bin/pip install --upgrade pip
#          venv/bin/pip install -r Tick_Data_Folder/requirements.txt
#
#  4.  Run auth_helper.py ONCE to generate the first access_token:
#          cd /home/ubuntu/Stock-Predictor/Tick_Data_Folder
#          ../venv/bin/python auth_helper.py --request-token <YOUR_TOKEN>
#
#  5.  Install this service file:
#          sudo cp /home/ubuntu/Stock-Predictor/Tick_Data_Folder/tick_harvester.service \
#                  /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable tick_harvester
#          sudo systemctl start tick_harvester
#
#  6.  Verify it is running:
#          sudo systemctl status tick_harvester
#          sudo journalctl -u tick_harvester -f
#          tail -f /home/ubuntu/Stock-Predictor/Tick_Data_Folder/harvester.log
#
#  7.  Stop gracefully (flushes CSV buffer before exit):
#          sudo systemctl stop tick_harvester
#
# ── DAILY ACCESS-TOKEN REFRESH (Sharekhan token expires at midnight IST) ─────
#
#  The TokenReminderThread emails you at 23:00 before each trading day.
#  After you receive the email, log in at the Sharekhan URL, copy the
#  RequestToken, then run on EC2 (or automate via cron):
#
#      cd /home/ubuntu/Stock-Predictor/Tick_Data_Folder
#      ../venv/bin/python auth_helper.py --request-token <REQUEST_TOKEN>
#      sudo systemctl restart tick_harvester
#
#  Cron automation (runs at 08:45 IST every weekday — token must be
#  pre-populated in config.json by that time):
#
#      crontab -e
#      # Append (time is IST because TZ=Asia/Kolkata is already set):
#      45 8 * * 1-5  sudo systemctl restart tick_harvester
#
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Sharekhan Live Tick Harvester
Documentation=https://github.com/Sharekhan-API/shareconnectpython

# Wait for full network + NTP time-sync before starting.
# time-sync.target ensures datetime.now() returns the correct IST time at boot.
After=network-online.target nss-lookup.target time-sync.target
Wants=network-online.target

[Service]
Type=simple

User=ubuntu
Group=ubuntu

# ── Paths ─────────────────────────────────────────────────────────────────────
WorkingDirectory=/home/ubuntu/Stock-Predictor/Tick_Data_Folder

ExecStart=/home/ubuntu/Stock-Predictor/venv/bin/python \
          /home/ubuntu/Stock-Predictor/Tick_Data_Folder/tick_harvester.py

# ── Timezone ──────────────────────────────────────────────────────────────────
# Force IST for all datetime.now() / date.today() calls inside the process.
# This is the single most important EC2 setting — without it the timestamps
# in the CSV files and holiday checks would be in UTC (off by 5 h 30 m).
Environment=TZ=Asia/Kolkata

# ── Restart policy ────────────────────────────────────────────────────────────
# on-failure: restart if the process crashes (non-zero exit).
# The holiday-sleep loop keeps the process alive overnight so systemd does
# NOT need to restart it during weekends / holidays.
# 'always' is intentionally avoided — it would mask expired-token errors.
Restart=on-failure
RestartSec=15

# ── Graceful shutdown ─────────────────────────────────────────────────────────
# SIGTERM triggers request_shutdown() → flushes the CSV buffer before exit.
# Give the process up to 30 s to finish flushing before sending SIGKILL.
KillSignal=SIGTERM
TimeoutStopSec=30

# ── Logging ───────────────────────────────────────────────────────────────────
# stdout/stderr → journald  (readable via journalctl)
# The script also writes to harvester.log directly for convenience.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tick_harvester

# ── Resource limits (t3.small: 2 GB RAM, 2 vCPU) ─────────────────────────────
# No memory cap — process is free to use available RAM.
CPUQuota=80%

[Install]
WantedBy=multi-user.target
